# .config/systemd/user/poweralertd.service.d/10-nostartevents.conf
[Service]
ExecStart=
ExecStart=/usr/bin/poweralertd -s
