# .config/systemd/user/sway.service.d/10-exit-target.conf
# https://gitlab.gnome.org/GNOME/gnome-session/-/commit/9de6e40f12e8878f524f8d429d85724c156a0517
[Service]
Slice=-.slice
